# php/crontab
*/1 * * * * /usr/local/php -q /var/www/commands/_cron.php >> /var/log/cron.log 2>&1
